# deploy/raw-editor.service
[Unit]
Description=RAW photo editor
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=joel
Group=joel
WorkingDirectory=/home/joel/dev/raw-editor
Environment=NODE_ENV=production
Environment=PATH=/home/joel/.nvm/versions/node/v24.13.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/joel/.nvm/versions/node/v24.13.0/bin/node server/app.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillSignal=SIGTERM
SyslogIdentifier=raw-editor

[Install]
WantedBy=multi-user.target
